# morientes.ini
[program:morientes]
command=/data/fiftyheis/morientes.py
redirect_stderr=true
stdout_logfile=/var/log/morientes.out
